# roles/monitoring/templates/prometheus.service.j2
[Unit]
Description=Prometheus
Documentation=https://prometheus.io/docs/
After=network.target

[Service]
Type=simple
User={{ monitoring_user }}
Group={{ monitoring_group }}
ExecStart=/usr/local/bin/prometheus \
  --config.file=/etc/prometheus/prometheus.yml \
  --storage.tsdb.path=/var/lib/prometheus \
  --storage.tsdb.retention.time={{ prometheus_retention_days }}d \
  --web.listen-address=0.0.0.0:{{ prometheus_listen_port }}

Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
